# rc.local/etc/rc.local
#!/bin/bash
sudo chown agung:agung /sys/class/backlight/amdgpu_bl0/*
exit 0
